# extra-files/lxc-dedicated-cpu-load/lxcfs.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/lxcfs -l /var/lib/lxcfs/